# webhook-handler.socket
[Socket]
ListenStream=/run/webhook-handler.sk

[Install]
WantedBy=sockets.target
